# logrotate/.config/systemd/user/logrotate-agent-logs.timer
[Unit]
Description=Weekly rotation of agent logs

[Timer]
OnCalendar=weekly
Persistent=true

[Install]
WantedBy=timers.target
