# logrotate/.config/systemd/user/logrotate-agent-logs.service
[Unit]
Description=Rotate agent logs (dippy audit, Claude Code hook approvals)

[Service]
Type=oneshot
ExecStart=/usr/sbin/logrotate --state %h/.local/state/logrotate-agent-logs.state %h/.config/logrotate/agent-logs.conf
